# knot/rootfs/etc/s6-overlay/s6-rc.d/create-sshd-host-keys/up
/etc/s6-overlay/scripts/create-sshd-host-keys
