# knot/rootfs/etc/s6-overlay/s6-rc.d/knotserver/run
#!/command/with-contenv ash

exec s6-setuidgid git /usr/bin/knot server
